# meta-hypercar/recipes-demo/linux-app/files/hypercar-app.service
[Unit]
Description=Start HyperCar POSIX app

[Service]
ExecStart=/bin/bash /usr/bin/hypercar-wrapper.sh

[Install]
WantedBy=multi-user.target